# deploy/devbox/systemd/hermes-telegram-deploy.timer
[Unit]
Description=Run hermes-telegram auto-deploy every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=15s

[Install]
WantedBy=timers.target
